# etc/php/crontab
0 0 * * * console event:playing-states:update --active
