# etc/php/crontab
10 0 * * * console steam:games:import
20 0 * * * console event:playing-states:update --active
